# crontab/readme
#运行用户：root

# 每天凌晨2：30执行数据库备份(前一天)
30 2 * * * /app/www/weitingba/data/crontab/backup_db.sh > /dev/null 2>&1
